# scripts/com.kemmer.df-pms-protel-adapter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- DF-PMS-Protel-Adapter macOS LaunchAgent [CRUX-MK]
     Welle-36 Foundation-DF (14. Foundation-DF). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kemmer.df-pms-protel-adapter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/make/Projects/dark-factories/df-pms-protel-adapter/scripts/run-df-pms-protel-adapter-mac.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/make/Projects/dark-factories/df-pms-protel-adapter/runs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/make/Projects/dark-factories/df-pms-protel-adapter/runs/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/make/Projects/dark-factories/df-pms-protel-adapter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DF_PMS_PROTEL_REAL_ENABLED</key>
        <string>false</string>
        <key>DF_PMS_PROTEL_TENANT_ID</key>
        <string>hildesheim</string>
        <key>DF_PMS_PROTEL_TIMEOUT_S</key>
        <string>30</string>
    </dict>
</dict>
</plist>
